# roles/dyndns/files/dy-update.timer
[Unit]
Description=Run dy.fi dynamic DNS update

[Timer]
OnCalendar=daily
Persistent=true
Unit=dy-update.service

[Install]
WantedBy=timers.target
